# scripts/com.polymarket.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/seanmurray/Projects/Poly_Up_Down/.venv/bin/python</string>
        <string>-m</string>
        <string>src.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seanmurray/Projects/Poly_Up_Down</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/seanmurray/Projects/Poly_Up_Down/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Restart automatically if the process exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 10s before restarting after a crash -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Stdout / stderr logs -->
    <key>StandardOutPath</key>
    <string>/Users/seanmurray/Projects/Poly_Up_Down/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seanmurray/Projects/Poly_Up_Down/logs/bot.error.log</string>
</dict>
</plist>
